# crontab
* * * * * /usr/bin/php /var/www/html/front/cron.php &>/dev/null
